# scripts/resources/discord_bot.conf
description "Minecraft Server"
start on runlevel[2345]
stop on runlevel[016]

chdir /home/username/discord-minecraft-server-manager/discord-bot/
exec "node main.js"
respawn
