# scripts/resources/minecraft.conf
description "Minecraft Server"
start on runlevel[2345]
stop on runlevel[016]

chdir /home/username/app/
exec "java -Xms2G -Xmx2G -jar /home/username/app/server.jar nogui"
respawn
